# etc/supervisord.d/statsd.ini
[program:statsd]
command = /usr/bin/statsd /opt/qnib/statsd/config.js
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
